# deploy/waybills.service
[Unit]
Description=Waybills — Kaspi picking service
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/opt/waybills
EnvironmentFile=/opt/waybills/.env
ExecStart=/opt/waybills/.venv/bin/uvicorn app.main:app --host 127.0.0.1 --port 8090 --proxy-headers --forwarded-allow-ips=127.0.0.1
Restart=on-failure
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
